# bacula-fd.init
#! /bin/sh
#
# bacula       This shell script takes care of starting and stopping
#	       the bacula File daemon.
#
# chkconfig: 2345 91 99
# description: It comes by night and sucks the vital essence from your computers.
#

# Source function library
. /etc/rc.d/init.d/functions

case "$1" in
    start)
       msg_starting "Bacula File daemon"
       daemon /usr/sbin/bacula-fd $2 -c /etc/bacula/bacula-fd.conf
       RETVAL=$?
       [ $RETVAL -eq 0 ] && touch /var/lock/subsys/bacula-fd
       ;;
    stop)
       msg_stopping "Bacula File daemon"
       killproc /usr/sbin/bacula-fd
       RETVAL=$?
       echo
       [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/bacula-fd
       ;;
    restart)
       $0 stop && sleep 5
       $0 start
       RETVAL=$?
       ;;
    status)
       status /usr/sbin/bacula-fd
       ;;
    *)
       echo "Usage: $0 {start|stop|restart|status}"
       exit 1
       ;;
esac
exit $RETVAL
